# crontab
*/5 * * * * /usr/local/bin/backup.sh 2>&1 > /tmp/backup.log
*/2 * * * * /usr/local/bin/mount_nas.sh
*/2 * * * * /usr/local/bin/sync_journal.sh
